# packages/wibed-battlemesh-experiment/files/etc/uci-defaults/battlemesh.sh
#!/bin/sh

/usr/bin/wbm-config > /tmp/wbm-config.log 2>&1
